# yardmaster.service
[Unit]
Description=Yardmaster task worker
After=network.target redis-server.service

[Service]
Type=simple
User=jared
WorkingDirectory=/home/jared/code/gibson-ops/yardmaster
ExecStart=/home/jared/.nvm/versions/node/v22.22.2/bin/node --import tsx/esm src/queue/worker-entry.ts
Restart=on-failure
RestartSec=10
Environment=NODE_ENV=production
Environment=NO_COLOR=1
Environment=NODE_OPTIONS=--max-old-space-size=4096 --heapsnapshot-near-heap-limit=3
Environment=PATH=/home/jared/.local/bin:/home/jared/.nvm/versions/node/v22.22.2/bin:/home/linuxbrew/.linuxbrew/bin:/usr/local/bin:/usr/bin:/bin
Environment=HOME=/home/jared

[Install]
WantedBy=multi-user.target
